# config/platform/com.hermass.weekend-repo-hygiene.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermass.weekend-repo-hygiene</string>
    <key>Program</key>
    <string>/Users/lv111101/Documents/hermass-observer-product/scripts/weekend_repo_hygiene_check.sh</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/lv111101/Documents/hermass-observer-product/logs/weekend_repo_hygiene_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/lv111101/Documents/hermass-observer-product/logs/weekend_repo_hygiene_stderr.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
